# crontab
sudo crontab -e
//exemplo
@reboot /usr/bin/node /home/tomas/ocaepic/updatebgp.js